# src/CodemixerApp/Resources/com.codecave.Codemixer.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Canonical LaunchAgent template for codemixerd.
  The GUI installer substitutes __CODEMIXERD_PATH__ at install time.
  Manual install:
    sed "s|__CODEMIXERD_PATH__|$(which codemixerd)|" \
      src/CodemixerApp/Resources/com.codecave.Codemixer.daemon.plist \
      > ~/Library/LaunchAgents/com.codecave.Codemixer.daemon.plist
    launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.codecave.Codemixer.daemon.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.codecave.Codemixer.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>__CODEMIXERD_PATH__</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/codemixerd.stdout</string>
    <key>StandardErrorPath</key>
    <string>/tmp/codemixerd.stderr</string>
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
